# files/supervisor/nslcd.conf
[program:nslcd]
command=/usr/sbin/nslcd -n -g
autostart=true
autorestart=true
