# files/supervisor/ntpd.conf
[program:ntpd]
command=/usr/sbin/ntpd -n -g
autostart=true
autorestart=true
