# files/supervisor/ntpd.conf
[program:ntpd]
command=/usr/sbin/ntpd
numprocs=1
autostart=true
autorestart=true
